# nova-finance/backend/supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/app/logs/supervisord.log
pidfile=/app/logs/supervisord.pid

[program:django]
command=gunicorn nova_backend.wsgi:application --bind 0.0.0.0:8000 --workers 4 --worker-class gevent --worker-connections 1000 --max-requests 1000 --max-requests-jitter 100 --timeout 30 --keep-alive 5 --preload
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/django.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:celery_worker]
command=celery -A nova_backend worker -l INFO -c 2
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/celery.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:celery_beat]
command=celery -A nova_backend beat -l INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/celery-beat.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface